# home/private_dot_config/profile.d/telemetry.sh

## Block ads and other analytics
## https://consoledonottrack.com/
export ADBLOCK=true
export DISABLE_OPENCOLLECTIVE=true
export DO_NOT_TRACK=1

# Disable Homebrew telemetry
# https://docs.brew.sh/Analytics
export HOMEBREW_NO_ANALYTICS=1

## Disable .NET SDK telemetry
## https://learn.microsoft.com/en-us/dotnet/core/tools/telemetry
export DOTNET_CLI_TELEMETRY_OPTOUT=1

## Disable Gatsby telemetry
## https://www.gatsbyjs.org/docs/telemetry/
export GATSBY_TELEMETRY_DISABLED=1

## Disable NextJs telemetry
## https://nextjs.org/telemetry
export NEXT_TELEMETRY_DISABLED=1

## Disable NextJs telemetry
## https://keystonejs.com/docs/reference/telemetry
export KEYSTONE_TELEMETRY_DISABLED=1

## Set ddev instrumentation_opt_in=false
export DDEV_NO_INSTRUMENTATION=true
